# bitrix-distr/crontab/root
# do daily/weekly/monthly maintenance
# min   hour    day     month   weekday command
*/15    *       *       *       *       run-parts /etc/periodic/15min
0       *       *       *       *       run-parts /etc/periodic/hourly
0       2       *       *       *       run-parts /etc/periodic/daily
0       3       *       *       6       run-parts /etc/periodic/weekly
0       5       1       *       *       run-parts /etc/periodic/monthly

# bitrix
# * * * * * test -f /var/www/www/bitrix/modules/main/tools/cron_events.php && { php -d sendmail_path="msmtp -t -i -a default" -f /var/www/www/bitrix/modules/main/tools/cron_events.php; } > /dev/stdout
* * * * * test -f /var/www/www/bitrix/php_interface/cron_events.php && { php -f /var/www/www/bitrix/php_interface/cron_events.php; } > /dev/stdout
